# systemd/wd-rac.service
[Unit]
Description=WsprDaemon Remote Access Channel (frpc tunnel)
Documentation=https://github.com/HamSCI/rac
After=network-online.target
Wants=network-online.target
# Stays inert until configured: frpc.toml only exists once the operator fills
# in the gw2 user/token/remotePort assignment from the WsprDaemon admin.
ConditionPathExists=/etc/sigmond/frpc.toml

[Service]
Type=simple
ExecStart=/usr/local/sbin/frpc -c /etc/sigmond/frpc.toml
Restart=always
RestartSec=30s
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
